# contrib/anywhereinput-server.service
[Unit]
Description=AnywhereInput Remote Control Server
Documentation=https://github.com/Z-Hussein/AnywhereInput
After=network.target

[Service]
Type=simple
User=%i
Group=%i
ExecStart=/usr/bin/python3 -m anywhereinput --tunnel local
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=anywhereinput

# Security hardening (non-breaking)
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/anywhereinput /home/%i/.local/share/anywhereinput
ProtectHome=read-only

[Install]
WantedBy=default.target
